# supervisor-msz-reverb.conf
; Copy to server: sudo cp /var/www/MSZ/supervisor-msz-reverb.conf /etc/supervisor/conf.d/msz-reverb.conf
; Then: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start msz-reverb

[program:msz-reverb]
process_name=%(program_name)s
command=php /var/www/MSZ/artisan reverb:start
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/MSZ/storage/logs/reverb.log
stopwaitsecs=3600
